# deploy/route53-frontend.service
# Route 53 Clone — Next.js frontend (production server)
# Install to: /etc/systemd/system/route53-frontend.service
#
# Runs `next start`, which serves the output of `npm run build`. Build must be
# run (with the right NEXT_PUBLIC_API_URL) BEFORE starting this service.

[Unit]
Description=Route 53 Clone Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/route53-clone/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
